# systemd/log-collector.timer
[Unit]
Description=Run log collector every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
